# docker/application/supervisor/caddy.ini
[program:caddy]
autorestart=true
autostart=true
command=caddy run --config /etc/caddy/Caddyfile --adapter caddyfile
directory=/var/www
startretries=3
stderr_logfile=/var/log/supervisor/caddy.err.log
stdout_logfile=/var/log/supervisor/caddy.out.log
user=root